# config/com.ollama.service.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ollama.service</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/ollama</string>
        <string>serve</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <!-- System Environment -->
        <key>HOME</key>
        <string>/Users/<OLLAMA_USER></string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <!-- Ollama Configuration -->
        <key>OLLAMA_HOST</key>
        <string>0.0.0.0:11434</string>
        <key>OLLAMA_MODELS</key>
        <string>/Users/<OLLAMA_USER>/.ollama/models</string>
        <key>OLLAMA_NUM_PARALLEL</key>
        <string>8</string>
        <key>OLLAMA_KEEP_ALIVE</key>
        <string>30m</string>
        <key>OLLAMA_FLASH_ATTENTION</key>
        <string>true</string>
        <key>OLLAMA_MAX_LOADED_MODELS</key>
        <string>4</string>
        <key>OLLAMA_NOPRUNE</key>
        <string>true</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <!-- Run as your user -->
    <key>UserName</key>
    <string><OLLAMA_USER></string>
    <key>GroupName</key>
    <string>staff</string>
    <key>StandardErrorPath</key>
    <string>/Users/<OLLAMA_USER>/mac-studio-server/logs/ollama.err</string>
    <key>StandardOutPath</key>
    <string>/Users/<OLLAMA_USER>/mac-studio-server/logs/ollama.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/<OLLAMA_USER>/.ollama</string>
</dict>
</plist> 